# smartdox-streamlit.service
[Unit]
Description=SmartDOX Streamlit Application
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/smartdox
Environment="PATH=/opt/smartdox/venv/bin"
Environment="API_BASE_URL=http://localhost:8000"
ExecStart=/opt/smartdox/venv/bin/streamlit run app.py \
    --server.port=8501 \
    --server.address=0.0.0.0 \
    --server.headless=true \
    --logger.level=info
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
